# target/linux/pistachio/base-files/etc/init.d/resetbootcount
#!/bin/sh /etc/rc.common

START=99

start() {
	# Mark successful boot by resetting boot count
	echo 0 > $(ls -1 /sys/devices/platform/*uboot-count/bootcount)
}
